# processor.conf
# Place in /etc/init

description "Running Video Converter"

start on net-device-up

respawn

exec java -jar /usr/local/bin/processor.jar --server.port=8080