# web.conf
# Place in /etc/init

description "Running Video Web Converter"

start on net-device-up

respawn

exec java -jar /usr/local/bin/web.jar --server.port=10000